# systemd/control-panel-horizon.service
[Unit]
Description=Liberu Control Panel Horizon Dashboard
After=network.target mysql.service redis.service

[Service]
Type=simple
User=www-data
Group=www-data
Restart=always
RestartSec=5
StartLimitInterval=0

# Working directory
WorkingDirectory=/var/www/control-panel

# Execute Horizon
ExecStart=/usr/bin/php /var/www/control-panel/artisan horizon

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/control-panel/storage

# Resource limits
LimitNOFILE=4096
MemoryMax=1G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=control-panel-horizon

[Install]
WantedBy=multi-user.target
